# pretzel-server/pretzel-server.service.example
# Copy to /etc/systemd/system/pretzel-server.service, adjust User and paths, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pretzel-server
#
# Listens on :3001 (see index.js). Start before remote-ui; tv-relay is independent.
# Optional LIFX HTTP API proxy: /lifx/* (set token below).

[Unit]
Description=Pretzel Pi API (TTS, volume, weather, LIFX proxy on :3001)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pretzel/pretzel-server
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
# Optional: match repo root VERSION after deploy (see README).
# Environment=PRETZEL_STACK_VERSION=1.5.5
# Optional: Settings page / curl admin — must match header from UI or client.
# Environment=PRETZEL_SETTINGS_PASSCODE=change_me
# Optional: git pull repo root (default: parent of this directory).
# Environment=PRETZEL_REPO_ROOT=/home/pi/pretzel
# Optional LIFX Cloud API proxy (/lifx/*):
# Environment=LIFX_API_TOKEN=your_token
# Environment=LIFX_API_URL=https://api.lifx.com/v1

[Install]
WantedBy=multi-user.target
