# tv-relay/tv-relay.service.example
# Copy to /etc/systemd/system/tv-relay.service, adjust User and paths, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tv-relay
#
# This is separate from pretzel-server (port 3001). Wyat NEXT_PUBLIC_TV_RELAY_URL must
# point at this service (port 3000 by default).

[Unit]
Description=Pretzel LG TV relay (HTTP :3000 to Wyat; WebSocket to TV)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=william
WorkingDirectory=/home/william/pretzel/tv-relay
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
# Optional: match repo root VERSION after deploy (see README).
# Environment=PRETZEL_STACK_VERSION=1.6.1
# After changing the registration manifest, re-pair: clear CLIENT_KEY below, add:
#   Environment=TV_FORCE_PAIRING=1
# then `sudo systemctl daemon-reload && sudo systemctl restart tv-relay`, accept the TV
# prompt, copy CLIENT_KEY from `journalctl -u tv-relay -n 30`, then set CLIENT_KEY=... and
# remove TV_FORCE_PAIRING. Fixes getPointerInputSocket 401 insufficient permissions.

[Install]
WantedBy=multi-user.target
