# remote-ui/remote-ui.service.example
# 1. Clone or pull the repo on the Pi, then install deps and build the UI (after each pull):
#      cd ~/pretzel && git pull
#      cd ~/pretzel/remote-ui && npm ci && npm run build
#    (`dist/` is required before `node server.cjs`; it is not committed. Use `npm install` only if you omit a lockfile.)
#
# 2. Copy this file to systemd and edit User + WorkingDirectory to match your Pi account and repo path:
#      sudo cp ~/pretzel/remote-ui/remote-ui.service.example /etc/systemd/system/remote-ui.service
#      sudo nano /etc/systemd/system/remote-ui.service
#
# 3. Enable and start:
#      sudo systemctl daemon-reload
#      sudo systemctl enable --now remote-ui
#      sudo systemctl status remote-ui
#
# Guest URL (same LAN): http://<pi-ip>:8080/
# Use a DHCP reservation or mDNS for a stable URL in QR codes.
# Start tv-relay (:3000) and pretzel-server (:3001) before this service (or add After= in [Unit]).

[Unit]
Description=Pretzel LAN remote UI (Vite React dist + proxy on 8080)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace `pi` and path with your login and where you cloned `pretzel`.
User=pi
WorkingDirectory=/home/pi/pretzel/remote-ui
ExecStart=/usr/bin/node server.cjs
Restart=always
RestartSec=5
Environment=NODE_ENV=production
# Optional: match repo root VERSION after deploy (see README).
# Environment=PRETZEL_STACK_VERSION=1.6.0
# Environment=REMOTE_UI_PORT=8080
# Environment=TV_RELAY_URL=http://127.0.0.1:3000
# Environment=PRETZEL_SERVER_URL=http://127.0.0.1:3001

[Install]
WantedBy=multi-user.target
